# discord-bot.service
[Unit]
Description=Discord Content Announcement Bot Service
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=xush
Environment="DISPLAY=:99"
Environment="NODE_ENV=production"
Environment="PATH=/home/xush/.nvm/versions/node/v22.17.0/bin:/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=/home/xush/discord-bot

# Kill any existing Xvfb and start new one, then start the bot
ExecStartPre=/bin/bash -c 'pkill Xvfb || true; sleep 2; Xvfb :99 -screen 0 1024x768x24 -ac +extension GLX +render -noreset & sleep 2'
ExecStart=/bin/bash -c 'cd /home/xush/discord-bot && source /home/xush/.nvm/nvm.sh && DISPLAY=:99 node index.js'
ExecStopPost=/bin/bash -c 'pkill Xvfb || true'

# Enhanced restart configuration
Restart=always
RestartSec=10s
RestartPreventExitStatus=0

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Logging configuration  
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bot

[Install]
WantedBy=multi-user.target